# etc/cron.d/incrediblepbx-firewall
*/10 5-22 * * * root /usr/local/sbin/ipchecker > /dev/null 2>&1
